# systemd/nc-fullscan.service
[Unit]
Description=Daily ClamAV full scan of NAS data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nc_fullscan
